# backup-nextcloud/files/backup_nextcloud_inc.service
[Unit]
Description=Execute incremental Nextcloud backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup_incremental
