# backup-nextcloud/files/backup_nextcloud_full.service
[Unit]
Description=Execute full Nextcloud backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup_full
