# backup-nextcloud/files/backup_nextcloud_inc.timer
[Unit]
Description=Run backup_nextcloud_inc.service every day at 2h05

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri,Sat *-*-* 02:05:00
